# release/src-ra/linux/vendors/Ralink/RT2880/inittab
::sysinit:/etc_ro/rcS
console::respawn:/bin/sh
